# resources/crontab
0 * * * * kr python3 -m kube_resource_report /output > /dev/null 2> /output/cronerror.log
